# telegram_bot_panelHOHO.service
[Unit]
Description=Telegram Bot Panel - Auto Start
After=network.target

[Service]
WorkingDirectory=/var/www/html
ExecStart=/usr/bin/python3 -m telegram_bot_panel
Restart=always
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

